# systemd/aideapulse-pipeline.timer
[Unit]
Description=Run AIdeaPulse pipeline nightly at 11pm Central

[Timer]
OnCalendar=*-*-* 23:00:00
TimezoneOfTimer=America/Chicago
Persistent=true

[Install]
WantedBy=timers.target
