# systemd/ideavault-pipeline.service
[Unit]
Description=IdeaVault Ingestion Pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/doubl/IdeaVault
ExecStart=/home/doubl/IdeaVault/pipeline/.venv/bin/python -m pipeline.main
EnvironmentFile=/home/doubl/IdeaVault/pipeline/.env
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
